# docs/com.echo.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.echo</string>
    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_WITH_YOUR_USER/echo</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/REPLACE_WITH_YOUR_USER/echo/.venv/bin/python</string>
        <string>-m</string>
        <string>app.main</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/echo.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/echo.err.log</string>
</dict>
</plist>
