# docs/echo.service.example
[Unit]
Description=Echo — Local Voice Transcription
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_YOUR_USER
WorkingDirectory=/path/to/echo
Environment="PATH=/path/to/echo/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/echo/.venv/bin/python -m app.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
